# nova-mixer.service
[Unit]
Description=nova-mixer — ChatMix for Arctis Nova Pro Wireless
After=pipewire.service
Wants=pipewire.service

[Service]
Type=simple
ExecStart=%h/.local/bin/nova-mixer
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
